# asterisk/systemd/fastagi@.service
[Service]
Type=simple
ExecStart=/opt/agi/bin/phpagi-fastagi.php
User=root
Group=root
StandardInput=socket
StandardOutput=socket
